# config/ntfy/com.openclaw.ntfy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd agent: keeps the self-hosted ntfy container running at login.

  Install:
    cp config/ntfy/com.openclaw.ntfy.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.openclaw.ntfy.plist

  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.openclaw.ntfy.plist
    rm ~/Library/LaunchAgents/com.openclaw.ntfy.plist

  Note: "restart: unless-stopped" in docker-compose.yml means Docker itself
  handles container-level restarts. This plist only ensures the compose stack
  is (re)started when the Mac boots / the user logs in.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.ntfy</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/docker</string>
        <string>compose</string>
        <string>--project-directory</string>
        <string>/Users/rreilly/.openclaw/workspace/config/ntfy</string>
        <string>up</string>
        <string>-d</string>
        <string>--remove-orphans</string>
    </array>

    <!-- Run 60 s after login to give Docker Desktop time to start -->
    <key>StartInterval</key>
    <integer>60</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/rreilly/.openclaw/logs/ntfy-launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/rreilly/.openclaw/logs/ntfy-launchd.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/rreilly</string>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
